# conky.service
# ~/.config/systemd/user/conky.service
[Unit]
Description=Conky System Monitor
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
# If you use a custom config file, you can specify it like this:
# ExecStart=/usr/bin/conky -c /path/to/your/conky.conf
ExecStart=/usr/bin/conky -c %h/.config/conky/thinkpad.conf
Restart=on-failure
RestartSec=3
# Optional: If Conky starts too fast before your desktop loads, uncomment the line below:
# ExecStartPre=/bin/sleep 5

[Install]
WantedBy=graphical-session.target
